# src/bin/server/deb/DEBIAN/postrm
#!/bin/bash

update-rc.d -f s3-server remove && exit 0







