# resources/dist/debian/synthese3-server/postrm
#!/bin/bash


update-rc.d -f synthese3 remove && exit 0







