# resources/dist/debian/synthese3-server/postrm
#!/bin/sh


update-rc.d synthese3-server remove






