# resources/dist/debian/synthese3-server/postinst
#!/bin/bash

update-rc.d -f synthese3 defaults

# Start server
test ! -x /etc/init.d/synthese3 || /etc/init.d/synthese3 start





