# resources/dist/debian/synthese3-server/postinst
#!/bin/sh

update-rc.d synthese3-server defaults

# Start server
test ! -x /etc/init.d/synthese3-server || /etc/init.d/synthese3-server start





